# packages/network/iw/udev.d/60-iw-regdomain.rules
KERNEL=="regulatory*", ACTION=="add", SUBSYSTEM=="platform", RUN+="/usr/lib/iw/setregdomain"
